# profile.d/clash_proxy_on.sh
source /etc/profile.d/clash.sh
proxy_on
echo ''
